# templates/etc/systemd/system/haproxy.service.d/50-after-sync-certificates.conf.j2
{{ ansible_managed | comment }}
[Unit]
Requires=sync-certificates.service
After=sync-certificates.service
