# templates/etc/systemd/system/sync-certificates.service.d/50-environment.conf.j2
{{ ansible_managed | comment }}
[Service]
Environment="SYNC_BUCKET={{ _acme.s3_bucket }}"
Environment="SYNC_PATH={{ _acme.ssl_path }}"
