# templates/etc/systemd/system/sync-certificates.service.d/51-reload-haproxy.conf.j2
{{ ansible_managed | comment }}
[Service]
ExecStopPost=/bin/sh -c '[ "$EXIT_STATUS" = "2" ] && /bin/systemctl --no-block try-reload-or-restart {{ _acme.instance_name }} || true'
